# scripts/podman/haoclaw.container.in
# Haoclaw gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into haoclaw's ~/.config/containers/systemd/
# {{HAOCLAW_HOME}} is replaced at install time.

[Unit]
Description=Haoclaw gateway (rootless Podman)

[Container]
Image=haoclaw:local
ContainerName=haoclaw
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{HAOCLAW_HOME}}/.haoclaw:/home/node/.haoclaw:Z
EnvironmentFile={{HAOCLAW_HOME}}/.haoclaw/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
