# scripts/systemd/haoclaw-auth-monitor.service
[Unit]
Description=Haoclaw Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/haoclaw/scripts/auth-monitor.sh
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=haoclaw-alerts

[Install]
WantedBy=default.target
